# boot/notengobattery-preboot/files/early_boot.init
#!/bin/sh /etc/rc.common
#
# Problems? Sugestions? Visit: https://notengobattery.com
#

NAME=early_boot
PROG=/usr/libexec/early_boot
START=00
STOP=99
USE_PROCD=1

start_service() {
  procd_open_instance early_boot
  procd_set_param command $PROG
  procd_append_param command boot
  procd_set_param stdout 0
  procd_set_param stderr 0
  procd_close_instance
}
